# systemd/hypebot.service
[Unit]
<<<<<<< HEAD
Description=Hyperliquid RSI Bot
=======
Description=HypeBot RSI Strategy Service
>>>>>>> 78efec8b1ea762142e6cb0b25bfea53bd9964f55
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/hypebot_rsiv3/main.py
WorkingDirectory=/opt/hypebot_rsiv3
Restart=always
User=root
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
